# server/deploy/systemd/memberroll-backup.timer
[Unit]
Description=Nightly memberroll backup

[Timer]
OnCalendar=*-*-* 03:17:00
Persistent=true

[Install]
WantedBy=timers.target
